# AWS/cicd/jenkins/jenkins-setup-ec2/ansible-playbook/playbooks/templates/jenkins-server.service.j2
[Unit]
Description=Jenkins Server

[Service]
User=jenkins
WorkingDirectory=/opt/jenkins/
ExecStart=/opt/jenkins/start-jenkins.sh
Type=simple
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
